# ops/vps/seo-agent-discord/seo-agent-repo-health.service
[Unit]
Description=SEO Agent repo checkout health check

[Service]
Type=oneshot
WorkingDirectory=/opt/ai-dashboard/apps/seo-agent-discord
ExecStart=/usr/bin/node /opt/ai-dashboard/apps/seo-agent-discord/repo-health-check.mjs
