# ops/vps/seo-runtime/seo-runtime.service
[Unit]
Description=SEO Runtime API
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/ai-dashboard/apps/seo-runtime
Environment=NODE_ENV=production
EnvironmentFile=/opt/ai-dashboard/apps/seo-runtime/.env
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
